# deployment/remote/ansible/templates/manfrod.service.j2
[Unit]
Description=Manfrod Phoenix App
After=paradedb.service
Requires=paradedb.service

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
WorkingDirectory={{ app_dir }}
ExecStart=/usr/bin/mix phx.server
Restart=on-failure
RestartSec=5

# Graceful shutdown: give Phoenix/Bandit time to drain connections and release the port
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Environment for graceful shutdown
Environment=RELEASE_DISTRIBUTION=none

[Install]
WantedBy=multi-user.target
