# deployment/remote/ansible/templates/paradedb.service.j2
[Unit]
Description=ParadeDB (Podman)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
ExecStartPre=-/usr/bin/podman rm -f {{ paradedb_container_name }}
ExecStart=/usr/bin/podman run --rm --name {{ paradedb_container_name }} \
  --env-file /etc/manfrod/paradedb.env \
  -p 127.0.0.1:{{ postgres_port }}:5432 \
  -v {{ paradedb_volume_name }}:/var/lib/postgresql \
  {{ paradedb_image }}
ExecStop=/usr/bin/podman stop -t 30 {{ paradedb_container_name }}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
